# tunnel_instance/etc/systemd/system/ssh.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=2022
